# serversetup/pdc/linuxclient/karoshi-3.0/configuration/etc/xdg/autostart/karoshi-detect-dvd.desktop
[Desktop Entry]
Type=Application
Name=Karoshi video DVD detection
Comment=Detect video DVDs and generate icon
OnlyShowIn=XFCE;
Exec=karoshi-detect-dvd
