# serversetup/pdc/linuxclient/karoshi-3.0/configuration/etc/xdg/autostart/lightsOn.desktop
[Desktop Entry]
Type=Application
Name=lightsOn
Comment=Delay screensaver when running fullscreen applications
OnlyShowIn=XFCE;
Exec=lightsOn.sh 120
