# serversetup/scripts/uptime/karoshi-uptime
#!/bin/bash
### BEGIN INIT INFO
# Provides:          karoshi-uptime
# Required-Start:    $local_fs $remote_fs $network $syslog
# Required-Stop:     $local_fs $remote_fs $network $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start karoshi-uptime
### END INIT INFO


#Set update interval in seconds
INTERVAL=5

PARAM=`echo $1 | tr -cd 'a-z'`

function show_help {
echo -e "\nMissing parameter\n\nkaroshi-uptime start\stop\status\n" 
}

function get_status {
#Get pid
UPTIMEPID=`pidof -x karoshi-uptime.sh`
if [ `echo $?` = 0 ]
then
	STATUSMSG="The karoshi-uptime daemon is running with process ID $UPTIMEPID."
	STATUS=running
else
	STATUSMSG="The karoshi-uptime daemon is stopped."
	STATUS=stopped
fi
}

if [ -z "$PARAM" ]
then
	show_help
	exit
fi

if [ $PARAM != start ] && [ $PARAM != stop ] && [ $PARAM != status ]
then
	show_help
	exit
fi

#Status
if [ $PARAM = status ]
then
	get_status
	echo $STATUSMSG
	exit
fi

#Stop
if [ $PARAM = stop ]
then
	get_status
	if [ $STATUS = running ]
	then
		killall karoshi-uptime.sh
	fi
	get_status
	echo $STATUSMSG
	exit
fi

#Start 
if [ $PARAM = start ]
then
	get_status
	if [ $STATUS = stopped ]
	then
		echo Starting daemon
		start-stop-daemon --start --background --exec /usr/sbin/karoshi-uptime.sh auto $INTERVAL &
		sleep 0.5
		get_status
	fi
	echo $STATUSMSG
	exit
fi

exit
