# serversetup/pdc/linuxclient/karoshi-2.2/utilities/karoshi-update.conf
# Karoshi Update
#
# This task pulls down new configuration files from a Karoshi server,
# to provide dynamic client updates.

description "Karoshi Configuration"
author      "Robin McCorkell <rmccorkell@karoshi.org.uk>"

start on (starting lightdm)

console log
task

pre-start script
	if [ ! -f /opt/karoshi/linuxclientsetup/variables ] || [ ! -f /opt/karoshi/linuxclientsetup/flags/domainset ]; then
		stop
		exit 0
	fi
end script

script
	. /opt/karoshi/linuxclientsetup/variables
	#Check if version has been set
	if [ -z "$LINUX_VERSION" ]; then
		echo Linux version not set
		LINUX_VERSION=notset
	fi

	#Make sure that netlogon share is not connected
	[ `mount | grep -c /tmp/netlogon` -gt 0 ] && umount /tmp/netlogon

	if [ ! -d /tmp/netlogon ]; then
		echo Creating netlogon folder
		mkdir /tmp/netlogon
	fi

	#Connect to the PDC to update the startup scripts
	echo Connecting to $PDC
	mount.cifs //"$PDC"/netlogon /tmp/netlogon -o guest,nounix 1>/dev/null

	#Update scripts from the server
	if [ -d /tmp/netlogon/linuxclient/$LINUX_VERSION/scripts/ ]; then
		cp -rf /tmp/netlogon/linuxclient/$LINUX_VERSION/scripts/* /opt/karoshi/linuxclientsetup/scripts/
		chmod 0755 /opt/karoshi/linuxclientsetup/scripts/*
	fi
	
	#Run client config
	[ -f /opt/karoshi/linuxclientsetup/scripts/client_config ] && /opt/karoshi/linuxclientsetup/scripts/client_config

	#Disconnect from $PDC
	echo Disconnecting from $PDC
	umount /tmp/netlogon
	
	#Restart winbind or nslcd, so they read the new configuration
	if [ "$DOMAINTYPE" = samba3 ]; then
		service winbind restart
	elif [ "$DOMAINTYPE" = samba4 ]; then
		service nslcd restart
	fi
end script

