# serversetup/all/reserveports/reserveports
#!/bin/bash
### BEGIN INIT INFO
# Provides:          reserveports
# Required-Start:    $local_fs $remote_fs $network $syslog
# Required-Stop:     $local_fs $remote_fs $network $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start reserveports
### END INIT INFO

for service in `ls -1 /etc/reserveports`
do
	#Attempt to reserve a port for each service listed
	for portdata in `cat /etc/reserveports/$service`
	do
		port=`echo "$portdata" | cut -d, -f1`
		options=`echo "$portdata" | cut -d, -f2`
		nc -"$options" "$port" </dev/null >/dev/null 2>/dev/null &
	done
done
exit
